# gw/etc/profile.d/history.sh
HISTSIZE=1000000
HISTFILESIZE=2000000
HISTTIMEFORMAT="%F %T "
